# etc/systemd/system/lowpan_monitor.service
[Unit]
Description=Create lowpan (IEEE802.15.4) monitor device
After=network.target
Conflicts=lowpan lowpan_route

[Service]
EnvironmentFile=/etc/default/lowpan_monitor
Type=oneshot
User=root
ExecStart=/usr/local/sbin/create_monitor $CHN
ExecStop=/usr/local/sbin/delete_monitor
RemainAfterExit=yes
